# paxd-libre
#!/sbin/openrc-run
  
depend() {
	need localmount
}
  
start() {
/usr/bin/paxd-libre
}
